# systemd/rpmlogger.timer
[Unit]
Description=Rpmlogger timer

[Timer]
OnCalendar=*-*-* 04:30:00
RandomizedDelaySec = 600

[Install]
WantedBy=multi-user.target